# etc/crontab
# battery-notify
*/1  *  *  *  * DISPLAY=:0.0 XAUTHORITY=~/.Xauthority ~/.bin/battery-notify.sh
