# etc/systemd/system/battery-notify.service
[Unit]
Description=Battery Notify

[Service]
Type=simple
ExecStart=/home/syui/.bin/battery-notify.sh
