# etc/systemd/system/battery-notify.timer
[Unit]
Description=Battery Notify

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
